# znc.ini

[program:znc]
user=
command=/usr/bin/znc -d /config/ -f
